# usr/lib/systemd/system/gstreamill.service.in
[Unit]
Description=The gstreamill encoder server

[Service]
Type=forking
PIDFile=/var/run/gstreamill/gstreamill.pid
ExecStartPre=$(bindir)/cleangstmill.sh
ExecStart=$(bindir)/gstreamill
ExecStop=/bin/kill -s SIGTERM $MAINPID
ExecStartPost=$(bindir)/startjobs.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
