# usr/lib/systemd/system/gstreamill.service
[Unit]
Description=The gstreamill encoder server

[Service]
Type=forking
PIDFile=/var/run/gstreamill/gstreamill.pid
ExecStart=/usr/local/bin/gstreamill
ExecStop=/bin/kill -s SIGUSR2 $MAINPID
ExecStartPost=/usr/local/bin/startjobs.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
