# etc/init/gstreamill.conf.in
description "gstreamill encoder server"
author "Zhang Ping <dqzhangp@163.com>"
 
start on (filesystem and net-device-up IFACE=lo)
stop on runlevel [!2345]
 
env DAEMON=$(bindir)/gstreamill
env PIDFILE=/var/run/gstreamill/gstreamill.pid

pre-start script
    rm -rf /var/run/gstreamill/gstreamill.pid
    mkdir -p /var/run/gstreamill/mq
    mount -t mqueue none /var/run/gstreamill/mq
    rm -rf /var/run/gstreamill/mq/*
    umount /var/run/gstreamill/mq
    rm -rf /dev/shm/*
end script

post-start script
    if [ -d /etc/gstreamill.d ]; then
        for job in $(ls /etc/gstreamill.d/jobs.d/*.job); do
            curl -H "Content-Type: application/json" --data @$job http://localhost:20118/admin/start
        done
    fi
end script
 
post-stop script
    killall -q gstreamill
end script

expect fork
respawn
respawn limit 10 5
 
exec $DAEMON
